# offregister_githook/conf/git_pull.upstart.conf
# THIS FILE IS NOT USED <YET>
description "%(DESCRIPTION)s"
author "%(AUTHOR)s"

start on (filesystem and net-device-up IFACE!=lo)
stop on runlevel [!2345]

env DAEMON="%(DAEMON)s"
env DAEMON_ARGS="%(DAEMON_ARGS)s"
env PID="%(PID)s"
env GIT_BRANCH="%(GIT_BRANCH)s"
env GIT_REMOTE="%(GIT_REMOTE)s"

exec $DAEMON $DAEMON_ARGS
script
        chdir $GIT_DIR
        exec git fetch
        exec git checkout -f $GIT_BRANCH
        exec git reset --hard $GIT_REMOTE/$GIT_BRANCH
end script
